# service/powerdns-confd/run
#!/bin/bash

. /bin/envd-load

confd -onetime -backend env -confdir `pwd`
sleep 600
