# service/powerdns/run
#!/bin/bash
sv start powerdns-confd || exit 1
exec pdns_server --daemon=no
